# supervisord.conf
[supervisord]
nodaemon=true

[program:python]
command=/usr/bin/python /preview-bot/run.py
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
